# wii-u-gc-adapter.service
[Unit]
Description=Wii U GameCube adapter daemon

[Service]
Type=exec
ExecStartPre=modprobe uinput
ExecStart=wii-u-gc-adapter
